# scripts/rc.local
#!/bin/sh -e

sh /opt/probe/start.sh

exit 0
