# scripts/rc.local
#!/bin/sh -e

#sh /opt/probe/start.sh
sh /opt/tindae-update/agent/bin/restart.sh
exit 0
