# scripts/rc.local
#!/bin/sh -e

/opt/wdog/watchdog -d
sh /opt/probe/start.sh


exit 0
